# setup-utils/silence.socket
[Unit]
Description=silence socket

[Socket]
ListenStream=/usr/share/silence/silence.sock

[Install]
WantedBy=sockets.target
